# systemd/motion_wide.service
[Unit]
Description=Pi Video System (Flask + sensors)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pi-video-system/current
ExecStart=/opt/pi-video-system/.venv/bin/python -m app.web
Restart=always
RestartSec=2
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
