# systemd/motion_wide.service
[Unit]
Description=Pi Video System (Flask + sensors)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pi-video-system/current
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=HOST=0.0.0.0
Environment=PORT=8080
Environment=USE_FLASK=1
ExecStart=/opt/pi-video-system/.venv/bin/python -m app.web
Restart=on-failure
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
